# 90-hdmi.rules
SUBSYSTEM=="drm", ACTION=="change", RUN+="/etc/acpi/hdmi-switch.sh"
